# systemd/pacifica-proxy-refresh.service
[Unit]
Description=Pacifica Proxy Refresh and Apply
After=network-online.target pacifica-ui-api.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/pacifica-flow
Environment=PACIFICA_PROXY_TEST_CONCURRENCY=96
Environment=PACIFICA_PROXY_TEST_TIMEOUT_MS=6000
Environment=PACIFICA_POSITIONS_PROXY_POOL_MIN_OK=200
Environment=PACIFICA_DETECTOR_PROXY_POOL_MIN_OK=200
ExecStartPre=/usr/bin/bash -n /root/pacifica-flow/scripts/refresh_live_proxy_pools.sh
ExecStart=/root/pacifica-flow/scripts/refresh_live_proxy_pools.sh
