# systemd/pacifica-wallet-indexers.target
[Unit]
Description=PacificaFlow Wallet Indexer Shards
PartOf=pacifica-flow.target

[Install]
WantedBy=pacifica-flow.target
