# systemd/pacifica-wallet-indexer.service
[Unit]
Description=PacificaFlow Wallet Indexer Worker
After=network-online.target
Wants=network-online.target
PartOf=pacifica-flow.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/pacifica-flow
EnvironmentFile=/root/pacifica-flow/config/runtime.env
Environment=PORT=3201
Environment=PACIFICA_SERVICE=api
Environment=PACIFICA_WS_ENABLED=false
Environment=PACIFICA_SNAPSHOT_ENABLED=false
Environment=PACIFICA_GLOBAL_KPI_ENABLED=false
Environment=PACIFICA_INDEXER_ENABLED=true
Environment=PACIFICA_INDEXER_RUNNER_ENABLED=true
Environment=PACIFICA_ONCHAIN_ENABLED=false
Environment=PACIFICA_ONCHAIN_RUNNER_ENABLED=false
Environment=PACIFICA_INDEXER_DISCOVERY_ONLY=false
Environment=PACIFICA_SOLANA_LOG_TRIGGER_ENABLED=false
Environment=PACIFICA_LIVE_TRADE_ATTRIBUTION_DB_ENABLED=false
Environment=PACIFICA_LIVE_TRADE_WALLET_ATTRIBUTION_ENABLED=false
Environment=PACIFICA_LIVE_TRADE_TX_ATTRIBUTION_ENABLED=false
Environment=PACIFICA_LIVE_WALLET_FIRST_ENABLED=false
Environment=PACIFICA_POSITION_LIFECYCLE_LOOP_ENABLED=false
Environment=PACIFICA_WALLET_METRICS_HISTORY_LOOP_ENABLED=false
Environment=NODE_OPTIONS=--max-old-space-size=4096
ExecStart=/usr/bin/env PATH=/usr/local/bin:/usr/bin:/bin PORT=3201 npm run worker:wallet-indexer
Restart=always
RestartSec=3
TimeoutStopSec=15

[Install]
WantedBy=pacifica-flow.target
