# systemd/pacifica-live-remote-pull.service
[Unit]
Description=PacificaFlow Live Remote Pull
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/pacifica-flow
ExecStart=/root/pacifica-flow/scripts/systemd/pull_live_remote.sh

