# systemd/pacifica-live-positions-census@.service
[Unit]
Description=PacificaFlow Live Positions Census Worker Shard %i
After=network-online.target
Wants=network-online.target
PartOf=pacifica-live-positions.target pacifica-flow.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/pacifica-flow
Nice=5
OOMScoreAdjust=0
ExecStart=/root/pacifica-flow/scripts/systemd/run_live_positions_census_shard.sh %i 4
Restart=always
RestartSec=3
TimeoutStopSec=30

[Install]
WantedBy=pacifica-live-positions.target
